# SteamOSRootReback-watcher.path
[Unit]
Description="Monitoring /etc/SteamOSRootReback/hooks.d/ for changes"

[Path]
PathChanged=/etc/SteamOSRootReback/hooks.d/
Unit=SteamOSRootReback-watcher.service

[Install]
WantedBy=multi-user.target
